# .ebextensions/crontab
0 * * * * docker restart $(docker ps -q)
