# .ebextensions/crontab
0,15,30,45 * * * * docker restart $(docker ps -q)
